# tool/运维/config/abbey.ini
[program:abbey_os]
command= /usr/jjhgame/game-server/abbey_os/bin/abbey_server
directory=/usr/jjhgame/game-server/abbey_os
redirect_stderr=true
stdout_logfile=/usr/jjhgame/game-server/log/abbey_os.log
stopasgroup=true

[program:wanba]
command= /usr/jjhgame/game-server/wanba/bin/wanba
directory=/usr/jjhgame/game-server/wanba/
redirect_stderr=true
stdout_logfile=/usr/jjhgame/game-server/log/wanba.log
stopasgroup=true

[program:browser]
command= /usr/jjhgame/game-server/browser/bin/browser
directory=/usr/jjhgame/game-server/browser/
redirect_stderr=true
stdout_logfile=/usr/jjhgame/game-server/log/browser.log
stopasgroup=true

[program:h5]
command= /usr/jjhgame/game-server/h5/bin/abbey_server
directory=/usr/jjhgame/game-server/h5/
redirect_stderr=true
stdout_logfile=/usr/jjhgame/game-server/log/h5.log
stopasgroup=true



[program:abbey_os]
command= /usr/jjhgame/game-server/vn88/bin/abbey_server
directory=/usr/jjhgame/game-server/vn88
redirect_stderr=true
stdout_logfile=/usr/jjhgame/game-server/log/vn88.log
stopasgroup=true